# supervisor/rsu-queue.conf
[program:rsu-queue]
; ── Queue worker สำหรับ RSU Health Platform ────────────────────────────────
; แก้ path ด้านล่างให้ตรงกับ server จริงก่อน deploy

; PHP binary: ปรับตาม PHP version ที่ Plesk ใช้
; ตัวอย่าง: /usr/bin/php8.2  หรือ  /opt/plesk/php/8.2/bin/php
command=/usr/bin/php8.2 /var/www/vhosts/YOUR_DOMAIN.com/httpdocs/artisan queue:work database \
        --sleep=3 \
        --tries=3 \
        --timeout=90 \
        --max-time=3600

; แก้ path + user ให้ตรงกับ Plesk subscription
directory=/var/www/vhosts/YOUR_DOMAIN.com/httpdocs
user=YOUR_PLESK_USER
numprocs=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true

; Log files (Plesk อาจใช้ /var/log/plesk/ หรือ path ใน subscription)
stdout_logfile=/var/www/vhosts/YOUR_DOMAIN.com/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/www/vhosts/YOUR_DOMAIN.com/logs/queue-worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

; Environment variables
environment=LARAVEL_SAIL=0
